# build/docker.conf
[supervisord]
nodaemon=true

[program:openfire]
command=/var/lib/openfire/start_openfire.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

